# src/shd/shd.service
[Unit]
Description=shd
After=network.target

[Service]
ExecStart=/usr/local/bin/shd
Type=simple

[Install]
WantedBy=default.target