# src/mathd/mathd.service
[Unit]
Description=mathd
After=network.target

[Service]
ExecStart=/usr/local/bin/mathd
Type=simple

[Install]
WantedBy=default.target
