# src/httpd/httpd.service
[Unit]
Description=httpd
After=network.target

[Service]
ExecStart=/usr/bin/httpd
Type=simple

[Install]
WantedBy=default.target