# src/noted/noted.service
[Unit]
Description=noted
After=network.target

[Service]
ExecStart=/usr/local/bin/noted
Type=simple

[Install]
WantedBy=default.target