# src/genericd/genericd.service
[Unit]
Description=genericd
After=network.target

[Service]
ExecStart=/usr/local/bin/genericd
Type=simple

[Install]
WantedBy=default.target