# src/fingerd/fingerd.service
[Unit]
Description=fingerd
After=network.target

[Service]
ExecStart=/usr/local/bin/fingerd
Type=simple

[Install]
WantedBy=default.target